# install/crontab/crontab.rules
* * * * * python /home/pi/code/IoTnode.py
* * * * * python /home/pi/code/IoTnodepir.py
